# dot_config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate Sway wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/sway/scripts/wallpaper-rotate.sh
